# bot.conf
[program:bot_test]
command=python -u ./nemoobot/nemoobot.py
directory=/home/nemoobot/ne-telegram-bot
stdout_logfile=/home/nemoobot/ne-telegram-bot/log/bot.log
redirect_stderr=true
